# home/aiden/.config/systemd/user/bspwm.target
[Unit]
Description=bspwm
Documentation=man:systemd.special(7)

BindsTo=xorg.target
PartOf=xorg.target

[Install]
WantedBy=xorg.target
